# deploy/toronto-viz.service
[Unit]
Description=Toronto TTC Visualization
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=/opt/toronto-viz
ExecStart=/opt/toronto-viz/deploy/run.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
